# configs/capture-agent.service
[Unit]
Description=capture-agent Network Packet Capture Daemon
Documentation=https://github.com/firestige/capture-agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Run as root for raw socket access (or use CAP_NET_RAW+CAP_NET_ADMIN)
User=root
Group=root

# Ensure data directories exist with correct permissions (ADR-031)
ExecStartPre=systemd-tmpfiles --create /etc/tmpfiles.d/capture-agent.conf
# Daemon mode (foreground for systemd)
ExecStart=/usr/local/bin/capture-agent daemon
ExecReload=/bin/kill -HUP $MAINPID

# Graceful shutdown (SIGTERM)
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30

# Auto-restart on failure
Restart=on-failure
RestartSec=5s

# Resource limits
LimitNOFILE=65536
LimitNPROC=512

# Security hardening (compatible with packet capture)
# Note: Some restrictions are relaxed due to raw socket requirements
NoNewPrivileges=false
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/var/lib/capture-agent /var/log/capture-agent /etc/capture-agent

# Required capabilities for packet capture
AmbientCapabilities=CAP_NET_RAW CAP_NET_ADMIN
CapabilityBoundingSet=CAP_NET_RAW CAP_NET_ADMIN CAP_DAC_OVERRIDE

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=capture-agent

# Working directory
WorkingDirectory=/var/lib/capture-agent

[Unit]
Description=Capture Agent Edge Packet Capture Agent
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/capture-agent
Restart=on-failure
RestartSec=5
User=root
Group=root
LimitNOFILE=65536
WorkingDirectory=/
EnvironmentFile=-/etc/default/capture-agent

[Install]
WantedBy=multi-user.target
