# configs/capture-agent.service
[Unit]
Description=capture-agent Network Packet Capture Daemon
Documentation=https://github.com/firestige/capture-agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Run as root for raw socket access (or use CAP_NET_RAW+CAP_NET_ADMIN)
User=root
Group=root

# Daemon mode (foreground for systemd)
ExecStart=/usr/local/bin/capture-agent daemon
ExecReload=/bin/kill -HUP $MAINPID

# Graceful shutdown (SIGTERM)
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30

# Auto-restart on failure
Restart=on-failure
RestartSec=5s

# Resource limits
LimitNOFILE=65536
LimitNPROC=512

# Directory management — systemd creates these BEFORE namespace setup,
# avoiding the NAMESPACE/226 failure caused by bind-mounting non-existent paths.
# StateDirectory       → /var/lib/capture-agent   (persistent state)
# LogsDirectory        → /var/log/capture-agent   (log files)
# ConfigurationDirectory → /etc/capture-agent     (config files, read-write overlay)
StateDirectory=capture-agent
LogsDirectory=capture-agent
ConfigurationDirectory=capture-agent

# Security hardening (compatible with packet capture)
# Note: Some restrictions are relaxed due to raw socket requirements
NoNewPrivileges=false
PrivateTmp=true
ProtectSystem=full
ProtectHome=true

# Required capabilities for packet capture
AmbientCapabilities=CAP_NET_RAW CAP_NET_ADMIN
CapabilityBoundingSet=CAP_NET_RAW CAP_NET_ADMIN CAP_DAC_OVERRIDE

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=capture-agent

# Working directory
WorkingDirectory=/var/lib/capture-agent

[Install]
WantedBy=multi-user.target

