# configs/otus.service
[Unit]
Description=Otus Network Packet Capture Daemon
Documentation=https://github.com/firestige/otus
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Run as root for raw socket access (or use CAP_NET_RAW+CAP_NET_ADMIN)
User=root
Group=root

# Ensure data directories exist with correct permissions (ADR-031)
ExecStartPre=systemd-tmpfiles --create /etc/tmpfiles.d/otus.conf
# Daemon mode (foreground for systemd)
ExecStart=/usr/local/bin/otus daemon
ExecReload=/bin/kill -HUP $MAINPID

# Graceful shutdown (SIGTERM)
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30

# Auto-restart on failure
Restart=on-failure
RestartSec=5s

# Resource limits
LimitNOFILE=65536
LimitNPROC=512

# Security hardening (compatible with packet capture)
# Note: Some restrictions are relaxed due to raw socket requirements
NoNewPrivileges=false
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/var/lib/otus /var/log/otus /etc/otus

# Required capabilities for packet capture
AmbientCapabilities=CAP_NET_RAW CAP_NET_ADMIN
CapabilityBoundingSet=CAP_NET_RAW CAP_NET_ADMIN CAP_DAC_OVERRIDE

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=otus

# Working directory
WorkingDirectory=/var/lib/otus

[Install]
WantedBy=multi-user.target
